# ellis/ellis.supervisord.conf
[program:ellis]
command=/etc/init.d/ellis run
autorestart=true
