# ellis/nginx.supervisord.conf
[program:nginx]
command=/etc/init.d/nginx start
startretries=20
