# ellis/nginx.supervisord.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=false
autorestart=true
