# homestead-prov/homestead-prov.supervisord.conf
[program:homestead-prov]
command=/etc/init.d/homestead-prov run
autostart=false
autorestart=true
