# ralf/chronos.supervisord.conf
[program:chronos]
command=/etc/init.d/chronos run
autorestart=true
