# ralf/ralf.supervisord.conf
[program:ralf]
command=/etc/init.d/ralf run
autorestart=true
