# aio/cassandra.supervisord.conf
[program:cassandra]
command=/etc/init.d/cassandra start
startretries=20
