# aio/ellis.supervisord.conf
[program:mysql]
command=/etc/init.d/mysql start
startretries=20

[program:ellis]
command=/etc/init.d/ellis start
startretries=20

[program:nginx]
command=/etc/init.d/nginx start
startretries=20
