# aio/homestead.supervisord.conf
[program:homestead]
command=/etc/init.d/homestead start
startretries=20

[program:homestead-prov]
command=/etc/init.d/homestead-prov start
startretries=20
