# aio/homer.supervisord.conf
[program:homer]
command=/etc/init.d/homer start
startretries=20
