# aio/sprout.supervisord.conf
[program:sprout]
command=/etc/init.d/sprout start
startretries=20

[program:memcached]
command=/etc/init.d/memcached start
startretries=20
