# astaire/astaire.supervisord.conf
[program:astaire]
command=/etc/init.d/astaire run
startretries=100
autostart=false
autorestart=true
