# astaire/memcached.supervisord.conf
[program:memcached]
command=/usr/bin/memcached -v -e 'ignore_vbucket=true;cache_size=512000000' -p 11211 -u memcache -c 4096
autostart=false
autorestart=true
