# chronos/chronos.supervisord.conf
[program:chronos]
command=start-stop-daemon --start --quiet --pidfile /var/run/chronos/chronos.pid --name chronos --exec /usr/bin/chronos -- --daemon
autostart=false
autorestart=true
