# chronos/chronos.supervisord.conf
[program:chronos]
command=/etc/init.d/chronos run
autostart=false
autorestart=true
