# chronos/snmpd.supervisord.conf
[program:snmpd]
command=/etc/init.d/snmpd start
autostart=false
autorestart=true
