# homer/cassandra.supervisord.conf
[program:cassandra]
command=/usr/sbin/cassandra -f
autostart=false
autorestart=true
