# homer/homer.supervisord.conf
[program:homer]
command=/etc/init.d/homer run
autorestart=true
