# homestead/cassandra.supervisord.conf
[program:cassandra]
# We've configured Cassandra to always run as a daemon, so we can issue "start"
# here and it won't fork
command=/etc/init.d/cassandra start
autorestart=true
