# homestead/cassandra.supervisord.conf
[program:cassandra]
command=/usr/sbin/cassandra -f
autorestart=true
