# homestead/homestead.supervisord.conf
[program:homestead]
command=/etc/init.d/homestead run
autostart=false
autorestart=true
# homestead kills itself if it can't talk to cassandra - give it a few retries
startretries=5000
