# base/alarm-agent.supervisord.conf
[program:alarm-agent]
command=/usr/share/clearwater/bin/clearwater-snmp-alarm-agent-wrapper
startsecs=0
