# base/monit.supervisord.conf
[program:monit]
command=/usr/bin/monit -I -c /etc/monit/monitrc
autostart=true
autorestart=true
startretries=10
