# base/snmpd.supervisord.conf
[program:snmpd]
command=/etc/init.d/snmpd start
startsecs=0
