# bono/bono.supervisord.conf
[program:bono]
command=/etc/init.d/bono run
autostart=false
autorestart=true
