# bono/bono.supervisord.conf
[program:bono]
command=/etc/init.d/bono run
autorestart=true
