# bono/bono.supervisord.conf
[program:bono]
command=/etc/init.d/bono start
startretries=20
