# bono/restund.supervisord.conf
[program:restund]
command=/etc/init.d/restund run
autorestart=true
