# bono/restund.supervisord.conf
[program:restund]
command=/etc/init.d/restund start
startretries=20
