# cassandra/cassandra.supervisord.conf
[program:cassandra]
command=/etc/init.d/cassandra start
autostart=false
# Cassandra won't start until the cluster manager runs - retry until that happens
startretries=5000
autorestart=true
