# cassandra/cassandra.supervisord.conf
[program:cassandra]
command=/usr/bin/start_cassandra.sh
autostart=false
# Cassandra won't start until the cluster manager runs - retry until that happens
startretries=5000
autorestart=true
