# sprout/sprout.supervisord.conf
[program:sprout]
command=/etc/init.d/sprout run
autostart=false
autorestart=true
